# puppet/netprobe/templates/netprobe_init.erb
#!/bin/bash
# chkconfig: 2345 80 20
# processname: netprobe_<%= @port %>
# description: Init script for netprobe_<%= @port %> \
# This script will start, stop, restart and return the status of a \
# netprobe instance using port <%= @port %>. Script was generated by \
# a puppet script. Please check the puppet master for more details.
#
#

# source function library
. /etc/rc.d/init.d/functions

RETVAL=0
prog="netprobe_<%= @port %>"
port=<%= @port %>
home="/opt/netprobe"
user="svc.itrs"

probestop() {
        PID=`ps -ef | grep "netprobe_$port.xml"| grep -v grep |awk '{ print $2 }'`
        echo "Stopping $prog: "
        if [ $UID -ne 0 ]; then
                RETVAL=1
                failure
        elif [[ -z $PID ]];then
                RETVAL=5
                echo "No netprobe on port $port running"
        else
                kill -9 $PID
                [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/${prog}
        fi;
        return $RETVAL
}


probestatus(){
        PID=`ps -ef | grep "netprobe_$port.xml"| grep -v grep | awk '{ print $2 }'`
        if [[ -z $PID ]];then
                echo "No netprobe on port $port running"
                RETVAL=5
        else
                echo "Service $prog running: pid $PID"
                RETVAL=0
        fi
}

case "$1" in
  start)
        echo -n "Starting $prog: "
        daemon --user=$user /opt/netprobe/start_netprobe $port
        touch /var/lock/subsys/$prog
        ;;
  stop)
        probestop
        ;;
  restart)
        probestop
	killall /opt/netprobe/netprobe.linux_64
        $0 start
        ;;
  status)
        probestatus
        ;;
  *)
        echo $"Usage: $0 {start|stop|status|restart}"
        RETVAL=1
esac

exit $RETVAL

